# linux-container/rootfs/home/sslab/raring-armhf-container/rootfs/var/lib/dpkg/info/opencryptoki.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pkcs_slot.1.gz
/usr/share/man/man1/pkcsconf.1.gz
/usr/share/man/man1/pkcs11_startup.1.gz
/usr/share/man/man5
/usr/share/man/man5/pk_config_data.5.gz
/usr/share/man/man8
/usr/share/man/man8/pkcsslotd.8.gz
/usr/share/man/man7
/usr/share/man/man7/opencryptoki.7.gz
/usr/share/doc
/usr/share/doc/opencryptoki
/usr/share/doc/opencryptoki/FAQ
/usr/share/doc/opencryptoki/AUTHORS
/usr/share/doc/opencryptoki/copyright
/usr/share/doc/opencryptoki/TODO.Debian
/usr/share/doc/opencryptoki/TODO
/usr/share/doc/opencryptoki/README.gz
/usr/lib
/usr/lib/opencryptoki
/usr/lib/opencryptoki/stdll
/usr/lib/pkcs11
/usr/sbin
/usr/sbin/pkcs11_startup
/usr/sbin/pkcs_slot
/usr/sbin/pkcsslotd
/usr/sbin/pkcsconf
/var
/var/lib
/var/lib/opencryptoki
/etc
/etc/init.d
/etc/init.d/opencryptoki
/usr/share/doc/opencryptoki/changelog.Debian.gz
/usr/lib/opencryptoki/stdll/PKCS11_SW.so
/usr/lib/opencryptoki/stdll/libpkcs11_tpm.so
/usr/lib/opencryptoki/stdll/libpkcs11_sw.so
/usr/lib/pkcs11/libopencryptoki.so
/usr/lib/pkcs11/PKCS11_API.so
/etc/pkcs11
