# linux-container/rootfs/home/sslab/raring-armhf-container/rootfs/var/lib/dpkg/info/trousers.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/tcsd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/tcsd.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/trousers
/usr/share/doc
/usr/share/doc/trousers
/usr/share/doc/trousers/AUTHORS
/usr/share/doc/trousers/README.selinux
/usr/share/doc/trousers/copyright
/usr/share/doc/trousers/examples
/usr/share/doc/trousers/examples/system.data.auth
/usr/share/doc/trousers/examples/system.data.noauth
/usr/share/doc/trousers/TODO.Debian
/usr/share/doc/trousers/TODO
/usr/share/doc/trousers/changelog.Debian.gz
/usr/share/doc/trousers/README.gz
/usr/share/doc/trousers/NICETOHAVES
/usr/sbin
/usr/sbin/tcsd
/usr/sbin/ps_convert
/usr/sbin/ps_inspect
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/45-trousers.rules
/var
/var/lib
/var/lib/tpm
/etc
/etc/init.d
/etc/init.d/trousers
/etc/tcsd.conf
