# linux-container/rootfs/var/lib/dpkg/info/lxc.list
/.
/etc
/etc/dnsmasq.d-available
/etc/dnsmasq.d-available/lxc
/etc/init
/etc/init/lxc-instance.conf
/etc/init/lxc-net.conf
/etc/init/lxc.conf
/etc/lxc
/etc/lxc/default.conf
/etc/lxc/auto
/etc/bash_completion.d
/etc/bash_completion.d/lxc.sh
/etc/apparmor.d
/etc/apparmor.d/abstractions
/etc/apparmor.d/abstractions/lxc
/etc/apparmor.d/abstractions/lxc/start-container
/etc/apparmor.d/abstractions/lxc/container-base
/etc/apparmor.d/lxc
/etc/apparmor.d/lxc/lxc-default
/etc/apparmor.d/lxc/lxc-default-with-nesting
/etc/apparmor.d/usr.bin.lxc-start
/etc/apparmor.d/lxc-containers
/etc/init.d
/etc/default
/etc/default/lxc
/usr
/usr/bin
/usr/bin/lxc-freeze
/usr/bin/lxc-execute
/usr/bin/lxc-unshare
/usr/bin/lxc-start
/usr/bin/lxc-create
/usr/bin/lxc-kill
/usr/bin/lxc-destroy
/usr/bin/lxc-checkpoint
/usr/bin/lxc-restart
/usr/bin/lxc-cgroup
/usr/bin/lxc-unfreeze
/usr/bin/lxc-stop
/usr/bin/lxc-aa-custom-profile
/usr/bin/lxc-wait
/usr/bin/lxc-start-ephemeral
/usr/bin/lxc-attach
/usr/bin/lxc-clone
/usr/bin/lxc-monitor
/usr/bin/lxc-netstat
/usr/bin/lxc-ps
/usr/bin/lxc-ls
/usr/bin/lxc-device
/usr/bin/lxc-info
/usr/bin/lxc-checkconfig
/usr/bin/lxc-console
/usr/bin/lxc-shutdown
/usr/bin/lxc-version
/usr/share
/usr/share/doc
/usr/share/doc/lxc
/usr/share/doc/lxc/examples
/usr/share/doc/lxc/examples/lxc-macvlan.conf
/usr/share/doc/lxc/examples/lxc-no-netns.conf
/usr/share/doc/lxc/examples/lxc-vlan.conf
/usr/share/doc/lxc/examples/lxc-empty-netns.conf
/usr/share/doc/lxc/examples/lxc-phys.conf
/usr/share/doc/lxc/examples/lxc-complex.conf
/usr/share/doc/lxc/examples/lxc-veth.conf
/usr/share/doc/lxc/copyright
/usr/share/lxc
/usr/share/lxc/hooks
/usr/share/lxc/hooks/mountcgroups
/usr/share/lxc/hooks/mountecryptfsroot
/usr/share/lxc/lxc.functions
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/lxc
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_lxc.py
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/lxc.7.gz
/usr/share/man/man5
/usr/share/man/man5/lxc.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/lxc-device.1.gz
/usr/share/man/man1/lxc-netstat.1.gz
/usr/share/man/man1/lxc-unfreeze.1.gz
/usr/share/man/man1/lxc-kill.1.gz
/usr/share/man/man1/lxc-ls.1.gz
/usr/share/man/man1/lxc-monitor.1.gz
/usr/share/man/man1/lxc-unshare.1.gz
/usr/share/man/man1/lxc-start.1.gz
/usr/share/man/man1/lxc-start-ephemeral.1.gz
/usr/share/man/man1/lxc-shutdown.1.gz
/usr/share/man/man1/lxc-cgroup.1.gz
/usr/share/man/man1/lxc-ps.1.gz
/usr/share/man/man1/lxc-restart.1.gz
/usr/share/man/man1/lxc-checkpoint.1.gz
/usr/share/man/man1/lxc-console.1.gz
/usr/share/man/man1/lxc-stop.1.gz
/usr/share/man/man1/lxc-attach.1.gz
/usr/share/man/man1/lxc-wait.1.gz
/usr/share/man/man1/lxc-info.1.gz
/usr/share/man/man1/lxc-freeze.1.gz
/usr/share/man/man1/lxc-clone.1.gz
/usr/share/man/man1/lxc-execute.1.gz
/usr/share/man/man1/lxc-version.1.gz
/usr/share/man/man1/lxc-create.1.gz
/usr/share/man/man1/lxc-checkconfig.1.gz
/usr/share/man/man1/lxc-destroy.1.gz
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/lxc
/usr/lib/arm-linux-gnueabihf/lxc/lxc-init
/var
/var/cache
/var/cache/lxc
/var/log
/var/log/lxc
/var/lib
/var/lib/lxc
/etc/init.d/lxc-instance
/etc/init.d/lxc
/etc/init.d/lxc-net
/usr/bin/lxc-halt
/usr/bin/lxc-list
/usr/share/doc/lxc/changelog.Debian.gz
