# linux-container/rootfs/var/lib/dpkg/info/opencryptoki.list
/var
/var/lib
/var/lib/opencryptoki
/etc
/etc/init.d
/etc/init.d/opencryptoki
