# linux-container/rootfs/var/lib/dpkg/info/trousers.list
/var
/var/lib
/var/lib/tpm
/etc
/etc/init.d
/etc/init.d/trousers
/etc/tcsd.conf
